# environment/coral/jailhouse/install/lib/systemd/system/uboot-imx.service
[Unit]
Description=U-Boot install service

[Service]
ExecStart=/usr/bin/install-uboot.sh

[Install]
WantedBy=multi-user.target