# environment/coral/jailhouse/install/lib/systemd/system/vitalsd.service
[Unit]
Description=vitalsd system vitality nurse
After=usb-gadget

[Service]
Exec=/usr/bin/vitalsd /dev/tty/ttyGS0

[Install]
WantedBy=basic.target
